# .config/systemd/user/indicator-kdeconnect.service
[Unit]
Description=KDE Connect daemon

[Service]
ExecStart=/usr/bin/indicator-kdeconnect

[Install]
WantedBy=default.target
