# .config/systemd/user/kdeconnectd.service
[Unit]
Description=KDE Connect daemon

[Service]
ExecStart=/usr/lib/kdeconnectd

[Install]
WantedBy=default.target
