# .config/systemd/user/mpDris2.service
[Unit]
Description=MPRIS V2.1 support for mpd 

[Service]
ExecStart=/usr/bin/mpDris2

[Install]
WantedBy=default.target
